# rockzero-automount.service
[Unit]
Description=RockZero Auto-Mount Service
Documentation=https://github.com/yourusername/rockzero
After=local-fs.target
Before=rockzero-service.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/rockzero/scripts/auto-mount-disks.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
NoNewPrivileges=no
ProtectSystem=full
ProtectHome=no

[Install]
WantedBy=multi-user.target
