# rockzero-safe-shutdown.service
[Unit]
Description=RockZero Safe Shutdown Service
Documentation=https://github.com/yourusername/rockzero
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Conflicts=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/rockzero/scripts/safe-shutdown.sh
TimeoutStartSec=90
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=yes
NoNewPrivileges=no

[Install]
WantedBy=shutdown.target reboot.target halt.target
