# rockzero-service.service
[Unit]
Description=RockZero Secure Service
After=network.target

[Service]
Type=simple
User=rockzero
Group=rockzero
WorkingDirectory=/opt/rockzero
Environment="RUST_LOG=info"
EnvironmentFile=/opt/rockzero/.env
ExecStart=/opt/rockzero/rockzero-service
Restart=always
RestartSec=10

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rockzero/data /opt/rockzero/uploads /DATA/AppData

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 硬件访问
SupplementaryGroups=video render plugdev

[Install]
WantedBy=multi-user.target
